# site-modules/profile/templates/netbox/netbox-rq.service.erb
[Unit]
Description=NetBox Request Queue Worker
Documentation=https://docs.netbox.dev/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

User=netbox
Group=netbox
WorkingDirectory=<%= @install_path %>

ExecStart=<%= @install_path %>/venv/bin/python3 <%= @install_path %>/netbox/manage.py rqworker high default low

Restart=on-failure
RestartSec=30
PrivateTmp=true

[Install]
WantedBy=multi-user.target
