# site-modules/profile/templates/cassandra/cassandra.service.erb
# Managed by puppet (class profile::cassandra); Changes will be lost.

[Unit]
Description=Cassandra
After=network.target

[Service]
User=<%= @cassandra_user %>
Group=<%= @cassandra_group %>

ExecStart=/opt/cassandra/bin/cassandra -f
LimitNOFILE=100000
LimitMEMLOCK=infinity
LimitNPROC=32768
LimitAS=infinity
Restart=on-failure
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target
