# site-modules/profile/templates/cassandra/cassandra.service.erb
# Managed by puppet (class profile::cassandra::node); Changes will be lost.

[Unit]
Description=Cassandra
After=network.target

[Service]
User=cassandra
Group=cassandra
Environment=JAVA_HOME=/usr/lib/jvm/java-8-openjdk-amd64
Environment=JVM_EXTRA_OPTS=-javaagent:<%= @exporter_path %>=@<%= @exporter_config %>
EnvironmentFile=-/etc/default/cassandra
ExecStart=/usr/sbin/cassandra -f
LimitNOFILE=100000
LimitMEMLOCK=infinity
LimitNPROC=32768
LimitAS=infinity
Restart=on-failure
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target
