# site-modules/profile/templates/sanoid/syncoid.service.erb
# Managed by puppet class profile::sanoid::syncoid_destination
# Changes will be overwritten

[Unit]
Description=ZFS dataset synchronization of <% @source %>
StartLimitIntervalSec=<%= @delay %>
StartLimitBurst=2

[Service]
Type=oneshot
User=root
Group=root
ExecStart=syncoid<%= @sync_options %> --sshkey <%= @ssh_key_filename %> root@<%= @source %> <%= @destination %>
Restart=on-failure
RestartSec=1s

[Install]
WantedBy=multi-user.target
