# site-modules/profile/templates/sanoid/syncoid.timer.erb
# Managed by puppet class profile::sanoid::syncoid_destination
# Changes will be overwritten
[Unit]
Description=ZFS dataset synchronization trigger

[Install]
WantedBy=timers.target

[Timer]
OnBootSec=1m
OnUnitActiveSec=<%= @frequency %>
RandomizedDelaySec=1min
FixedRandomDelay=true
Unit=<%= @service_basename %>.service
