# site-modules/profile/templates/swh/deploy/graph/swh-graph-shm-mount.service.erb
# Managed by puppet class profile::swh::deploy::graph
# Changes will be overwritten

[Unit]
Description=swh-graph RAM data cache in /dev/shm

[Service]
Type=oneshot
User=<%= @user %>
Group=<%= @group %>
RemainAfterExit=yes
ExecStart=mkdir -p /dev/shm/swh-graph/default
ExecStart=sh -c "ln -s /srv/softwareheritage/graph/latest/compressed/* /dev/shm/swh-graph/default"
ExecStart=sh -c "cp --remove-destination /srv/softwareheritage/graph/latest/compressed/graph.graph /dev/shm/swh-graph/default"
ExecStart=sh -c "cp --remove-destination /srv/softwareheritage/graph/latest/compressed/graph-transposed.graph /dev/shm/swh-graph/default"
ExecStop=rm -rf /dev/shm/swh-graph/default

[Install]
WantedBy=multi-user.target
