# site-modules/profile/templates/swh/deploy/graph/swh-graph.service.erb
# Managed by puppet class profile::swh::deploy::graph
# Changes will be overwritten

[Unit]
Description=swh-graph HTTP server
After=swh-graph-shm-mount.service
Requires=swh-graph-shm-mount.service

[Service]
<%- if !@sentry_dsn.nil? and !@sentry_dsn.empty? -%>
Environment=SWH_SENTRY_DSN=<%= @sentry_dsn %>
<%- end -%>
<%- if !@sentry_environment.nil? and !@sentry_environment.empty? -%>
Environment=SWH_SENTRY_ENVIRONMENT=<%= @sentry_environment %>
<%- end -%>
<%- if !@sentry_swh_package.nil? and !@sentry_swh_package.empty? -%>
Environment=SWH_MAIN_PACKAGE=<%= @sentry_swh_package %>
<%- end -%>
Type=simple
User=<%= @user %>
Group=<%= @group %>
ExecStart=/opt/swhgraph_venv/bin/swh graph -C <%= @config_file %> rpc-serve -g /dev/shm/swh-graph/default/graph

[Install]
WantedBy=multi-user.target
