# site-modules/profile/templates/swh/deploy/worker/swh-worker@.service.erb
# Managed by puppet - modifications will be overwritten
# In defined class profile::swh::deploy::worker::base

[Unit]
Description=Software Heritage Worker (%i)
After=network.target

[Service]
User=swhworker
Group=swhworker

Type=simple
# Following environment variables can be overriden in the respective
# swh-worker@<service>.service.d/parameters.conf
Environment=SWH_CONFIG_FILENAME=/etc/softwareheritage/%i.yml
Environment=SWH_LOG_TARGET=journal
Environment=CONCURRENCY=10
Environment=MAX_TASKS_PER_CHILD=5
Environment=LOGLEVEL=info
Environment=CELERY_HOSTNAME=<%= @celery_hostname %>
Environment=CELERY_WORKER_EXTRA_ARGS=
ExecStart=/usr/bin/python3 -m celery worker -n %i@${CELERY_HOSTNAME} --app=swh.scheduler.celery_backend.config.app --pool=prefork ${CELERY_WORKER_EXTRA_ARGS} --concurrency=${CONCURRENCY} --maxtasksperchild=${MAX_TASKS_PER_CHILD} -Ofair --loglevel=${LOGLEVEL} --without-gossip --without-mingle --without-heartbeat

KillMode=process
KillSignal=SIGTERM
TimeoutStopSec=15m

OOMPolicy=kill

Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
