# site-modules/profile/templates/phabricator/phabricator-phd.service.erb
# Phabricator Daemon unit file
# Managed by puppet class profile::phabricator
# Changes will be overwritten

[Unit]
Description=Phabricator Daemons
After=network.target mariadb.service apache2.service

[Service]
User=<%= @install_user %>
Group=<%= @install_user %>
Type=forking
Environment="PATH=/sbin:/usr/sbin:/usr/local/sbin:/usr/local/bin:/usr/bin:/bin"
ExecStart=<%= @phabricator_basepath -%>/phabricator/bin/phd start
ExecStop=<%= @phabricator_basepath -%>/phabricator/bin/phd stop
Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target
