# site-modules/profile/templates/swh/deploy/webapp/swh-webapp-update-savecodenow-statuses.service.erb
# Webapp Update SaveCodeNow status unit file
# Managed by puppet class profile::swh::deploy::webapp
# Changes will be overwritten

[Unit]
Description=Software Heritage Webapp Update SaveCodeNow Statuses
Requires=network.target
After=network.target

[Service]
User=<%= @user %>
Group=<%= @group %>
Environment=DJANGO_SETTINGS_MODULE=<%= @webapp_settings_module %>
Type=simple
ExecStart=/usr/bin/django-admin refresh_savecodenow_statuses

[Install]
WantedBy=multi-user.target
