# site-modules/profile/templates/swh/deploy/webapp/swh-webapp-update-savecodenow-statuses.timer.erb
[Install]
WantedBy=timers.target

[Timer]
OnCalendar=*:0/1
Unit=swh-webapp-update-savecodenow-statuses.service

[Unit]
Description=Recurring webapp update savecodenow statuses
