# site-modules/profile/templates/swh/deploy/webapp/update-savecodenow-statuses.timer.erb
[Install]
WantedBy=timers.target

[Timer]
OnCalendar=*:0/1
RandomizedDelaySec=1m
FixedRandomDelay=true
Unit=<%= @service_basename %>.service

[Unit]
Description=Recurring webapp update savecodenow statuses
