# site-modules/profile/templates/swh/deploy/scrubber/swh-scrubber-checker-postgres@.service.erb
# Scrubber checker postgres template unit file
# Managed by puppet class profile::swh::deploy::scrubber::checker::postgres
# Changes will be overwritten

[Unit]
Description=Software Heritage Scrubber Checker Postgres %i
After=network.target

[Service]
<%- if !@sentry_dsn.nil? and !@sentry_dsn.empty? -%>
Environment=SWH_SENTRY_DSN=<%= @sentry_dsn %>
<%- end -%>
<%- if !@sentry_environment.nil? and !@sentry_environment.empty? -%>
Environment=SWH_SENTRY_ENVIRONMENT=<%= @sentry_environment %>
<%- end -%>
<%- if !@sentry_swh_package.nil? and !@sentry_swh_package.empty? -%>
Environment=SWH_MAIN_PACKAGE=<%= @sentry_swh_package %>
<%- end -%>
Environment=SWH_SCRUBBER_CLI_EXTRA_ARGS="--help"
Environment=SWH_CONFIG_FILENAME=
User=<%= @user %>
Group=<%= @group %>
Type=simple
ExecStart=/usr/bin/swh scrubber check storage $SWH_SCRUBBER_CLI_EXTRA_ARGS

[Install]
WantedBy=multi-user.target
