# site-modules/profile/templates/swh/deploy/scheduler/swh-scheduler-update-metrics.timer.erb
[Install]
WantedBy=timers.target

[Timer]
OnCalendar=daily
Unit=swh-scheduler-update-metrics.service

[Unit]
Description=regular scheduler update-metrics refresh
