# site-modules/profile/templates/swh/deploy/scheduler/swh-scheduler-runner.service.erb
# Scheduler Runner unit file
# Managed by puppet class profile::swh::deploy::scheduler
# Changes will be overwritten

[Unit]
Description=Software Heritage scheduler runner
Requires=network.target rabbitmq-server.service
After=network.target rabbitmq-server.service

[Service]
User=<%= @user %>
Group=<%= @group %>
Type=simple
ExecStart=/usr/bin/swh-scheduler --config-file <%= @config_file %> --log-level <%= @runner_log_level %> runner --period 10
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
