# site-modules/profile/templates/swh/deploy/scheduler/swh-scheduler-update-metrics.timer.erb
[Install]
WantedBy=timers.target

[Timer]
Unit=swh-scheduler-update-metrics.service
OnCalendar=00/4:00
RandomizedDelaySec=1h
FixedRandomDelay=true

[Unit]
Description=regular scheduler update-metrics refresh
