# site-modules/profile/templates/swh/deploy/scheduler/swh-scheduler-listener.service.erb
# Scheduler Listener unit file
# Managed by puppet class profile::swh::deploy::scheduler
# Changes will be overwritten

[Unit]
Description=Software Heritage scheduler listener
Requires=network.target rabbitmq-server.service
After=network.target rabbitmq-server.service

[Service]
User=<%= @user %>
Group=<%= @group %>
Type=simple
Environment=SWH_LOG_TARGET=journal
ExecStart=/usr/bin/swh --log-level <%= @listener_log_level %> scheduler --config-file <%= @config_file %> listener
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
