# site-modules/profile/templates/maven_index_exporter/maven_index_exporter@.timer.erb
[Install]
WantedBy=timers.target

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=24h

[Unit]
Description=Run '<%= @template_name %>@.service' daily
