# templates/phabricator/phabricator-phd.service.erb
# Phabricator Daemon unit file
# Managed by puppet class profile::phabricator
# Changes will be overwritten

[Unit]
Description=Phabricator Daemons
After=network.target
Before=apache2.service

[Service]
User=<%= @phabricator_user %>
Group=<%= @phabricator_user %>
Type=oneshot
Environment="PATH=/sbin:/usr/sbin:/usr/local/sbin:/usr/local/bin:/usr/bin:/bin"
ExecStart=<%= @phabricator_basepath -%>/phabricator/bin/phd start
ExecStart=<%= @phabricator_basepath -%>/phabricator/bin/phd stop
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
