# site-modules/profile/files/swh/deploy/worker/swh-worker.service
[Unit]
Description=Software Heritage Worker (meta-task)

[Service]
Type=oneshot
ExecStart=/bin/true
ExecReload=/bin/true
RemainAfterExit=on

[Install]
WantedBy=multi-user.target
