# site-modules/profile/files/swh/deploy/worker/system-swh\x2dworker.slice
[Unit]
Description=SWH workers

[Slice]
MemoryMax=90%
