# site-modules/profile/files/swh/deploy/scrubber/swh-scrubber.slice
[Unit]
Description=SWH Scrubber services

[Slice]
MemoryMax=75%
