# templates/swh/deploy/journal/swh-journal-publisher.service.erb
# Journal Publisher unit file
# Managed by puppet class profile::swh::deploy::journal_publisher
# Changes will be overwritten

[Unit]
Description=Software Heritage Journal Publisher
Requires=network.target kafka.service
After=network.target kafka.service

[Service]
User=<%= @user %>
Group=<%= @group %>
Type=simple
ExecStart=/usr/bin/python3 -m swh.journal.publisher
Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target
