# site-modules/profile/templates/jenkins/agent/jenkins-agent.service.erb
[Unit]
Description=Jenkins agent
Documentation=https://wiki.jenkins.io/display/JENKINS/Distributed+builds
After=network.target

[Service]
EnvironmentFile=/etc/default/jenkins-agent
ExecStart=/usr/bin/java -cp <%= @jenkins_agent_jar %> hudson.remoting.jnlp.Main -headless -url $JENKINS_URL -workDir $AGENT_WORKDIR $JNLP_SECRET $AGENT_NAME

User=jenkins
Group=jenkins

KillMode=process

Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target
