# templates/swh/deploy/scheduler/ssh-ghtorrent.service.erb
# SSH Ghtorrent unit file
# Managed by puppet class profile::swh::deploy::scheduler::updater::consumer
# Changes will be overwritten

[Unit]
Description=SSH connection to GHTorrent infra
Requires=network-online.target
After=network-online.target

[Service]
User=<%= @consumer_user %>
Group=<%= @consumer_group %>
Type=simple
ExecStart=/usr/bin/ssh -i <%= @ghtorrent_private_key %> -L <%= @local_port %>:streamer.ghtorrent.org:5672 ghtorrent@streamer.ghtorrent.org
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
