# templates/swh/deploy/scheduler/swh-scheduler-updater-writer.service.erb
# Scheduler Updater Writer unit file
# Managed by puppet class profile::swh::deploy::scheduler::updater::writer
# Changes will be overwritten

[Unit]
Description=Software Heritage Scheduler Updater Writer
Requires=network.target 
After=network.target 

[Service]
User=<%= @writer_user %>
Group=<%= @writer_group %>
Type=oneshot
ExecStart=echo /usr/bin/python3 -m swh.scheduler.updater.writer

[Install]
WantedBy=multi-user.target
