# templates/swh/deploy/scheduler/swh-scheduler-updater-consumer-ghtorrent.service.erb
# Scheduler Updater Consumer unit file
# Managed by puppet class profile::swh::deploy::scheduler::updater::consumer
# Changes will be overwritten

[Unit]
Description=Software Heritage Scheduler Updater Consumer
Requires=network.target ssh-ghtorrent.service
After=network.target ssh-ghtorrent.service

[Service]
User=<%= @backend_user %>
Group=<%= @backend_group %>
Type=simple
ExecStart=/usr/bin/python3 -m swh.scheduler.updater.ghtorrent.cli
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
