# site-modules/profile/files/maven_index_exporter/maven_index_exporter.slice
[Unit]
Description=Maven index exporters

[Slice]
MemoryMax=75%
