# site-modules/profile/templates/swh/deploy/loader_opam/opam-manage-shared-state.service.erb
# Opam manage shared state unit file
# Managed by puppet class profile::swh::deploy::worker::loader_opam
# Changes will be overwritten

[Unit]
Description=Software Heritage Manage OPAM shared state
Requires=network.target
Before=swh-worker@lister.service swh-worker@loader_opam.service
After=network.target

[Service]
User=<%= @user %>
Group=<%= @group %>
Type=oneshot
ExecStart=<%= @opam_manage_state_script %>

[Install]
WantedBy=swh-worker@lister.service swh-worker@loader_opam.service

