# site-modules/profile/templates/swh/deploy/loader_opam/opam-manage-shared-state.timer.erb
[Install]
WantedBy=timers.target

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=24h
FixedRandomDelay=true
Unit=<%= @opam_manage_service_name %>.service

[Unit]
Description=Software Heritage Manage OPAM shared state timer
