# site-modules/profile/templates/swh/deploy/journal/swh-indexer-journal-client.service.erb
# Indexer Journal Client unit file
# Managed by puppet class profile::swh::deploy::indexer_journal_client
# Changes will be overwritten

[Unit]
Description=Software Heritage Indexer Journal Client
After=network.target

[Service]
User=<%= @user %>
Group=<%= @group %>
Type=simple
ExecStart=/usr/bin/swh indexer --config-file <%= @config_path %> journal-client
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
