# site-modules/profile/templates/swh/deploy/journal/swh-indexer-journal-client.service.erb
# Indexer Journal Client unit file
# Managed by puppet class profile::swh::deploy::indexer_journal_client
# Changes will be overwritten

[Unit]
Description=Software Heritage Indexer Journal Client
After=network.target

[Service]
<%- if !@sentry_dsn.nil? and !@sentry_dsn.empty? -%>
Environment=SWH_SENTRY_DSN=<%= @sentry_dsn %>
<%- end -%>
<%- if !@sentry_environment.nil? and !@sentry_environment.empty? -%>
Environment=SWH_SENTRY_ENVIRONMENT=<%= @sentry_environment %>
<%- end -%>
<%- if !@sentry_swh_package.nil? and !@sentry_swh_package.empty? -%>
Environment=SWH_MAIN_PACKAGE=<%= @sentry_swh_package %>
<%- end -%>
User=<%= @user %>
Group=<%= @group %>
Type=simple
ExecStart=/usr/bin/swh indexer --config-file <%= @config_path %> journal-client indexer origin-intrinsic-metadata
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
